# cron.conf
[program:cron]
user=sysadmin
command=/usr/sbin/cron -f
